# stage2/10-bt-headless-a2dp/files/pulseaudio.service
[Unit]
Description=Pulseaudio System Daemon
After=bluetooth.service

[Service]
Type=simple
ExecStart=/usr/bin/pulseaudio
PIDFile=/var/run/pulse/pid

[Install]
WantedBy=multi-user.target
